# templates/upstart_geodns.conf
description "geodns"

start on runlevel [2345]
stop on runlevel [!2345]

script
  exec /opt/geodns/geodns -config=/opt/geodns/dns
end script