# templates/upstart.conf
description "BukGet Web API"
author "steve@chigeek.com"

start on stopped rc RUNLEVEL=[345]
stop on runlevel [0156]

exec /usr/bin/bukget

respawn